# files/etc/periodic/daily/backup
#!/bin/bash

python /backup-cleaner.py -r 1 -v

ruby /backup_script.rb run
